# debian/zookeeper.cron.daily
#!/bin/sh

. /usr/share/zookeeper/bin/zkPurgeTxnLog.sh
